# xx.conf
[program:xx]
command=python3 /home/curator2/xx/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/xx/test.err.log
stdout_logfile=/home/curator2/xx/test.out.log